# podkop/files/etc/cron.d/podkop_auto_update
# Podkop automatic subscription updates
# This file is managed by podkop - do not edit manually

# Auto-update subscriptions every hour (will check internal timers)
0 * * * * root /usr/bin/podkop auto_update_subscriptions >/dev/null 2>&1
